# admin/gunicorn_supervisord.conf
[program:gunicorn_script]
command=/home/seckin/gunicorn_script.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/gunicorn_script.err.log
stdout_logfile=/var/log/gunicorn_script.out.log
